# deploy/systemd/opencrm-api.service
[Unit]
Description=OpenCRM Backend API (Elysia)
After=network.target postgresql.service redis-server.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/.openclaw/workspace/opencrm-app/apps/backend
EnvironmentFile=/home/ubuntu/.openclaw/workspace/opencrm-app/.env
Environment=APP_MODE=api
Environment=PATH=/home/ubuntu/.bun/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/env -u APP_MODE APP_MODE=api /home/ubuntu/.bun/bin/bun --no-env-file run dist/index.js
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/opencrm-api.log
StandardError=append:/var/log/opencrm-api.log

[Install]
WantedBy=multi-user.target
